# mywebapp.service
[Unit]
Description=mywebapp - Simple Inventory Service
After=network.target postgresql.service

[Service]
User=app
Group=app
WorkingDirectory=/opt/mywebapp

# Run database migration BEFORE starting the app
ExecStartPre=/opt/mywebapp/venv/bin/python3 /opt/mywebapp/migrate.py

# Start the Flask application
ExecStart=/opt/mywebapp/venv/bin/python3 /opt/mywebapp/app.py

Restart=always

[Install]
WantedBy=multi-user.target
